# templates/default/kafka.service.erb
[Unit]
Description=Kafka publish-subscribe messaging system
After=network.target

[Service]
User=<%= @node['confluent-platform']['kafka']['user'] %>
Group=<%= @node['confluent-platform']['kafka']['user'] %>
SyslogIdentifier=kafka
Restart=on-failure
ExecStart=/usr/bin/java \
<% @node['confluent-platform']['kafka']['jvm_opts'].merge(
  @node['confluent-platform']['kafka']['jmx_opts']).each do |key, value| %>
  <%= key %><% unless value.nil? %>=<%= value %><% end %> \
<% end %>
  -Dlog4j.configuration=file:/etc/kafka/log4j.properties \
  -cp /usr/share/java/kafka/* \
  kafka.Kafka \
  /etc/kafka/server.properties

[Install]
WantedBy=multi-user.target
