# templates/default/kafka-rest.service.erb
[Unit]
Description=The Kafka REST Proxy provides a RESTful interface to a Kafka
After=network.target

[Service]
User=<%= @node['confluent-platform']['rest']['user'] %>
Group=<%= @node['confluent-platform']['rest']['user'] %>
SyslogIdentifier=kafka-rest
Restart=on-failure
<% cp = [ 'confluent-common','rest-utils', 'kafka-rest' ].map do |path| %>
  <% "/usr/share/java/#{path}/*" %>
<% end %>
ExecStart=/usr/bin/java \
  <%= @node['confluent-platform']['rest']['heap_opts'] %> \
  <%= @node['confluent-platform']['rest']['performance_opts'] %> \
  <%= @node['confluent-platform']['rest']['jmx_opts'] %> \
  <%= @jmx_port %> \
  <%= @node['confluent-platform']['rest']['extra_opts'] %> \
  -Dlog4j.configuration=file:/etc/kafka-rest/log4j.properties \
  -cp <%= cp.join(':') %> \
  io.confluent.kafkarest.Main \
  /etc/kafka-rest/kafka-rest.properties

[Install]
WantedBy=multi-user.target
